# misc/rotate_logs.timer
[Unit]
Description=runs rotate_logs every day

[Timer]
OnCalendar=daily
Unit=rotate_logs.service

[Install]
WantedBy=multi-user.target
